# deploy/cloud/hermes-webui.service
# systemd unit for the cloud-hosted Hermes WebUI.
#
# Install: sudo cp deploy/cloud/hermes-webui.service /etc/systemd/system/
#          sudo systemctl daemon-reload
#          sudo systemctl enable --now hermes-webui
#
# Logs:    sudo journalctl -u hermes-webui -f
# Restart: sudo systemctl restart hermes-webui
#
# Edit the Environment= lines below before installing — at minimum
# HERMES_WEBUI_AUTH_MODE must be set, and the working directory must
# point at your hermes-installer checkout.

[Unit]
Description=Hermes WebUI (cloud-hosted at chat.neowow.studio)
After=network-online.target
Wants=network-online.target

[Service]
Type=exec

# ── Identity ──────────────────────────────────────────────────────────
# Recommend a dedicated unprivileged user for isolation. Create with:
#   sudo useradd -r -m -d /opt/hermes -s /bin/bash hermes
User=hermes
Group=hermes

# ── Working directory ────────────────────────────────────────────────
# Point this at the hermes-installer checkout. The user `hermes` must
# own this directory.
WorkingDirectory=/opt/hermes/hermes-installer
ExecStart=/bin/bash webui/start.sh

# ── Environment ──────────────────────────────────────────────────────
# Auth mode — neodomain delegates login to app.neowow.studio's OAuth
# flow. Cross-subdomain `neoToken` cookie carries the JWT.
Environment=HERMES_WEBUI_AUTH_MODE=neodomain

# Port the WebUI listens on (loopback only; Caddy proxies HTTPS to it).
Environment=HERMES_WEBUI_PORT=7891
Environment=HERMES_WEBUI_HOST=127.0.0.1

# Foreground mode tells bootstrap.py to NOT detach — systemd needs the
# main process to stay running.
Environment=HERMES_WEBUI_FOREGROUND=1

# State / config dir. Default ~/.hermes/webui inherits HOME from the
# `hermes` user (= /opt/hermes/.hermes). Override here if you want a
# different layout (e.g. shared NFS for HA).
# Environment=HERMES_WEBUI_STATE_DIR=/var/lib/hermes/webui

# UTF-8 everywhere (Chinese filenames in Workspace, etc.).
Environment=PYTHONIOENCODING=utf-8
Environment=PYTHONUTF8=1

# ── Resource limits ───────────────────────────────────────────────────
# Adjust to your ECS instance specs. The settings below assume a 2 GB
# RAM minimum; bump for heavier workloads.
MemoryMax=4G
CPUQuota=200%

# ── Restart policy ────────────────────────────────────────────────────
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=5

# Long-running connections (SSE, websockets) — give the process room to
# shut down gracefully on `systemctl stop`.
TimeoutStopSec=15
KillMode=mixed
KillSignal=SIGTERM

# ── Hardening ─────────────────────────────────────────────────────────
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/hermes
# ReadWritePaths=/var/lib/hermes  # uncomment if you set HERMES_WEBUI_STATE_DIR

[Install]
WantedBy=multi-user.target
